# package/webface/files/etc/crontabs/root
# syncs time
15	* * * *	/etc/init.d/timesync start >/dev/null 2>/dev/null
